# pokey.service
[Unit]
Description=Pokey Service

[Service]
ExecStart=/usr/local/bin/pokey.py
